# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf
user=postgres
autorestart=true
priority=1

[program:qdrant]
command=/usr/local/bin/qdrant --config-path /etc/qdrant/config.yaml
user=root
autorestart=true
priority=2

[program:ollama]
command=/usr/local/bin/ollama serve
user=root
autorestart=true
priority=3

[program:backend]
command=java -jar /app/backend/app.jar
user=root
autorestart=true
priority=10
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.err

[program:agent]
command=uvicorn main:app --host 0.0.0.0 --port 9000
directory=/app/agent
user=root
autorestart=true
priority=11
stdout_logfile=/var/log/agent.log
stderr_logfile=/var/log/agent.err

[program:nginx]
command=nginx -g 'daemon off;'
user=root
autorestart=true
priority=20
